# net-print/p910nd/files/init
#!/sbin/openrc-run

command=/usr/sbin/p910nd
command_args="${OPTIONS} ${NUMBER:-0}"
pidfile="/run/p910${NUMBER:-0}d.pid"

depent() {
	need net
}
